# roles/workstation/templates/deepcool-ak400.service

[Unit]
Description=Displaying the current temperature and utilisation on the Deepcool AK400

[Service]
Type=simple
User=root
ExecStart= /usr/bin/python3 /usr/bin/deepcool-ak400.py
Restart=always

[Install]
WantedBy=default.target